# tpa-alpha-bot/deploy/systemd-frontend.service
[Unit]
Description=TPA Alpha Frontend (Dev Server)
After=network.target

[Service]
Type=simple
User=ubuntu
WorkingDirectory=/home/ubuntu/apps/tpa-alpha-bot/frontend
Environment="PATH=/home/ubuntu/.nvm/versions/node/v18.0.0/bin:/usr/bin"
Environment="NODE_ENV=production"
ExecStart=/home/ubuntu/.nvm/versions/node/v18.0.0/bin/npm run preview
Restart=always
RestartSec=10

# For Production, use nginx to serve build output instead
# Build: npm run build
# Output: dist/
# Serve via nginx on port 80

[Install]
WantedBy=multi-user.target
