# tpa-alpha-bot.service
[Unit]
Description=TPA Alpha Discord Bot
After=network.target tpa-backend.service

[Service]
Type=simple
User=ubuntu
WorkingDirectory=/var/www/tpa-alpha/tpa-alpha-bot
Environment="PATH=/var/www/tpa-alpha/tpa-alpha-bot/backend/venv/bin"
EnvironmentFile=/var/www/tpa-alpha/tpa-alpha-bot/backend/.env
ExecStart=/var/www/tpa-alpha/tpa-alpha-bot/backend/venv/bin/python bot/bot.py
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
